# deploy/vertel-bot.service
[Unit]
Description=VerTel Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vertel
Group=vertel
WorkingDirectory=/opt/vertel-bot
ExecStart=/opt/vertel-bot/vertel_bot
Restart=always
RestartSec=3
Environment=VERTEL_INJECT_SAMPLE_START=0
Environment=VERTEL_TELEGRAM_LONG_POLL_TIMEOUT_SECONDS=25
Environment=VERTEL_TELEGRAM_REQUEST_TIMEOUT_SECONDS=35
Environment=VERTEL_POLL_MAX_ATTEMPTS=5
Environment=VERTEL_POLL_INITIAL_BACKOFF_MS=250
Environment=VERTEL_LOOP_SLEEP_MS=50
Environment=VERTEL_RATE_LIMIT_CAPACITY=5
Environment=VERTEL_RATE_LIMIT_REFILL_TOKENS=5
Environment=VERTEL_RATE_LIMIT_REFILL_SECONDS=10
EnvironmentFile=/etc/vertel-bot/vertel-bot.env

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/vertel-bot
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
